# files/docker-system-prune.timer
[Unit]
Description=Cleanup unused images and stopped docker containers
[Timer]
OnCalendar=weekly
[Install]
WantedBy=timers.target
